# pkg-build/scripts/postinstall
#!/bin/bash
# LFG postinstall — install LaunchAgents for the real user, symlink CLI
set -euo pipefail

# Determine the actual logged-in user (not root)
REAL_USER=$(stat -f "%Su" /dev/console 2>/dev/null || echo "${SUDO_USER:-$USER}")
if [[ -z "$REAL_USER" || "$REAL_USER" == "root" ]]; then
    # Fallback: first human user in /Users
    REAL_USER=$(ls /Users | grep -v "^Shared$" | grep -v "^\." | head -1)
fi

REAL_HOME=$(dscl . -read "/Users/$REAL_USER" NFSHomeDirectory 2>/dev/null | awk '{print $2}')
if [[ -z "$REAL_HOME" ]]; then
    REAL_HOME="/Users/$REAL_USER"
fi

echo "LFG postinstall: configuring for user $REAL_USER ($REAL_HOME)"

# --- LaunchAgents ---
LAUNCH_AGENTS="$REAL_HOME/Library/LaunchAgents"
mkdir -p "$LAUNCH_AGENTS"

for plist in /Users/Shared/lfg/io.lfg.*.plist; do
    [[ -f "$plist" ]] || continue
    BASENAME="$(basename "$plist")"
    # Patch all developer-specific paths for the installing user
    sed \
        -e "s|/Users/jeremiah/tools/@yj/lfg|/Users/Shared/lfg|g" \
        -e "s|/Users/jeremiah/.config/lfg|$REAL_HOME/.config/lfg|g" \
        -e "s|<string>/Users/jeremiah</string>|<string>$REAL_HOME</string>|g" \
        "$plist" > "$LAUNCH_AGENTS/$BASENAME"
    chown "$REAL_USER" "$LAUNCH_AGENTS/$BASENAME" 2>/dev/null || true
    # Load via launchctl bootstrap (modern macOS 13+) with fallback
    USER_UID=$(id -u "$REAL_USER" 2>/dev/null || echo "")
    if [[ -n "$USER_UID" ]]; then
        launchctl bootstrap "gui/$USER_UID" "$LAUNCH_AGENTS/$BASENAME" 2>/dev/null || \
            sudo -u "$REAL_USER" launchctl load "$LAUNCH_AGENTS/$BASENAME" 2>/dev/null || true
    else
        sudo -u "$REAL_USER" launchctl load "$LAUNCH_AGENTS/$BASENAME" 2>/dev/null || true
    fi
done

# --- CLI symlink ---
mkdir -p /usr/local/bin
ln -sf /Users/Shared/lfg/lfg /usr/local/bin/lfg 2>/dev/null || true

# --- DevDrive setup (best effort) ---
sudo -u "$REAL_USER" /Users/Shared/lfg/lfg devdrive setup 2>/dev/null || true

echo "LFG installed successfully for $REAL_USER."
exit 0
