# pkg-build/root/Users/Shared/lfg/io.lfg.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.lfg.helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jeremiah/tools/@yj/lfg/LFG Helper.app/Contents/MacOS/LFG Helper</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart unless the user explicitly quit (exit 0) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/jeremiah/.config/lfg/helper.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jeremiah/.config/lfg/helper.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jeremiah</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
