# pkg-build/root/Users/Shared/lfg/io.lfg.inbox-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.lfg.inbox-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeremiah/tools/@yj/lfg/scripts/inbox-watcher.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/jeremiah/.config/lfg/inbox/watcher-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jeremiah/.config/lfg/inbox/watcher-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jeremiah</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LFG_DIR</key>
        <string>/Users/jeremiah/tools/@yj/lfg</string>
        <key>LFG_INBOX_POLL</key>
        <string>5</string>
    </dict>
</dict>
</plist>
